# supervisor/privoxy.svc.conf
[program:privoxy]
directory=/etc/service/privoxy/
command=/etc/service/privoxy/run
redirect_stderr=true

# vi: ft=cfg

