# supervisor/tor.svc.conf
[program:tor]
directory=/etc/service/tor/
command=/etc/service/tor/run
redirect_stderr=true

# vi: ft=cfg

